# dotfiles/systemd/.config/systemd/user/cleanup-downloads.timer
[Unit]
Description=Daily cleanup of old Downloads

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
